# crontab
* * * * * cd /app && /usr/local/bin/python ./main.py > /proc/1/fd/1 2>/proc/1/fd/2
